# modules/base/.config/profile.d/base-env.sh
#
# Common environment user session
# NOTE: This works for bash, zsh, and other POSIX shells.
#

# Colors
export CLICOLOR=1

# Pager
export LESS="-F -g -i -M -R -S -w -X -z-4"
export LESSCHARSET="utf-8"
export PAGER="less"
