# mod/docker/.config/profile.d/docker.sh
#
# Docker profile
#

# Xdg directores
export DOCKER_CONFIG="${XDG_CONFIG_HOME}/docker"
