# modules/macos/Library/LaunchAgents/dotfiles.launchctl.env.plist
<!-- https://www.manpagez.com/man/5/launchd.plist/ -->
<?xml version="1.0" encoding="UTF-8" ?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dotfiles-launchctl-env</string>
  <key>EnableGlobbing</key>
  <true />
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>exec $HOME/.local/bin/launchctl-env.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true />
  <key>KeepAlive</key>
  <false />
  <key>StandardOutPath</key>
  <string>/tmp/dotfiles.launchctl.env.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/dotfiles.launchctl.env.log</string>
</dict>
</plist>
