# stow/.config/profile.d/stow.sh
# --------------------------------------------------------------
# Stow dotfiles helpers
# --------------------------------------------------------------

alias stowh="stow -t ~"
alias stowo="stow -t ~ -d opt"
alias ustow="stow -D"
alias ustowh="stow -D -t ~"
alias ustowo="stow -D -t ~ -d opt"
