# stow/.config/profile.d/stow.sh
# --------------------------------------------------------------
# Stow dotfiles helpers
# --------------------------------------------------------------

alias ustow="stow -D"
alias stowh="stow -t ~"
alias ustowh="stow -D -t ~"
alias stowp="stow -t ~ -d pkg"
alias ustowp="stow -D -t ~ -d pkg"
