# modules/vim/.config/profile.d/vim.sh
#
# Vim environment
#

# Remap ~/.vimrc for XDG specification
export VIMINIT=":source $XDG_CONFIG_HOME/nvim/init.vim"
