# modules/macos/.config/profile.d/macos.sh
#
# Initialize macOS environment
#

# Misc XDG specifications
export CP_HOME_DIR="${XDG_DATA_HOME}/cocoapods"
