# modules/homebrew/.config/profile.d/homebrew.sh
# Source the homebrew shellenv
source_file "${XDG_CONFIG_HOME:-$HOME/.config}/homebrew/shellenv.sh"
