# opt/gnupg/.config/profile.d/gpg.sh
#
# GnuPG
#

# XDG specification
export GNUPGHOME="$XDG_DATA_HOME"/gnupg
