# modules/_base/.config/profile.d/base-env.sh
#
# Common environment user session
#

# Editor
# export EDITOR="nvim"
# export VISUAL="nvim -b"

# Colors
export CLICOLOR=1

# Pager
export LESS="-F -g -i -M -R -S -w -X -z-4"
export LESSCHARSET="utf-8"
export PAGER="less"
