# idea/.config/profile.d/idea.sh
#
# Environment for JetBrains IDEA
#

# Change location of ~/idea.properties file to XDG spec
export IDEA_PROPERTIES="${XDG_CONFIG_HOME}/idea/idea.properties"
